# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:agent]
command=python -m backend.agent start
directory=/app
autostart=true
autorestart=true
startretries=10
stopwaitsecs=30
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:ui_server]
command=/bin/sh -c 'uvicorn backend.main:app --host 0.0.0.0 --port ${PORT:-8000} --workers 1'
directory=/app
autostart=true
autorestart=true
startretries=5
stopwaitsecs=15
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
